# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash
set -e

exec su -c crond -f

# EOF